# src/brain_cli/data/services/brain-dream.service
[Unit]
Description=Brain Knowledge Graph Dream Maintenance

[Service]
Type=oneshot
ExecStart=$BRAIN_PATH dream --phased
WorkingDirectory=$BRAIN_DIR/..

[Install]
WantedBy=default.target
